# systemd/sys_cleanup.service
[Unit]
Description=System cleanup script
Wants=sys_cleanup.timer

[Service]
Type=oneshot
ExecStart=/root/scripts/sys_cleanup.sh
