# systemd/system_update.service
[Unit]
Description=System Update Script

[Service]
Type=oneshot
ExecStart=/root/scripts/system_update.sh
